# .config/local/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/project/.config/local/uwsgi_socket.ini

[program:nginx]
command:nginx